# .dev/trainer.service
[Unit]
Description=G.O.D Trainer Service
After=network.target docker.service
Wants=docker.service
Requires=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/G.O.D
Environment="PATH=/root/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="ENV=DEV"
# Allow enough time for image downloads before startup.
TimeoutStartSec=30min
# Deploy log shipping before starting the server
ExecStartPre=/bin/bash -c 'source /root/.venv/bin/activate && cd /root/G.O.D && task deploy-trainer-logs'
# Pull trainer images before starting.
ExecStartPre=/bin/bash -c 'set -euo pipefail; docker_cli=$(command -v docker); [ -n "${docker_cli:-}" ] || { echo "god-trainer: docker not found in PATH" >&2; exit 1; }; for img in gradientsio/trainer-uploader:latest gradientsio/trainer-downloader:latest gradientsio/trainer-cacher-cleaner:latest gradientsio/model-prep:latest; do "$docker_cli" pull "$img"; done'
# Start uvicorn server directly
ExecStart=/bin/bash -c 'source /root/.venv/bin/activate && cd /root/G.O.D && uvicorn trainer.asgi:app --host 0.0.0.0 --port 8001 --env-file .trainer.env --log-level debug'
# Restart on failure with exponential backoff
Restart=on-failure
RestartSec=5
StartLimitInterval=300
StartLimitBurst=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=god-trainer

[Install]
WantedBy=multi-user.target

